# telldusd.service
[Unit]
Description=Telldus daemon

[Service]
Type=daemon
ExecStart=/usr/sbin/telldusd
Restart=always

[Install]
WantedBy=multi-user.target
